# wrt/MainHub/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

smbd -D
nmbd -D

echo 2 > /proc/irq/49/smp_affinity
echo 2 > /proc/irq/50/smp_affinity

exit 0
