# joe/http-bash-requests.timer
[Unit]
Description=Run HTTP Server 2 Minutes After Boot

[Timer]
OnBootSec=2min
Unit=http-bash-requests.service

[Install]
WantedBy=timers.target
