# wrt/MainHub/etc/crontabs/root
0 */6 * * * sleep "$(awk 'BEGIN{srand();print int(rand()*6*60)}')m" ; /root/checkspeed.sh
*/5 * * * * /root/domain.sh
#*/5 * * * * /root/host-update.sh
